# installer/pi-image/stage-openavc/01-install-openavc/files/openavc-firstboot.service
[Unit]
Description=OpenAVC First Boot Setup
After=network-online.target
ConditionPathExists=/var/lib/openavc/.firstboot

[Service]
Type=oneshot
ExecStart=/opt/openavc/scripts/openavc-firstboot.sh
RemainAfterExit=false

[Install]
WantedBy=multi-user.target
