# installer/pi-image/stage-openavc/01-install-openavc/files/openavc.service
[Unit]
Description=OpenAVC Room Control Server
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=openavc
Group=openavc
WorkingDirectory=/opt/openavc
ExecStartPre=-+/opt/openavc/update-helper.sh /var/lib/openavc
ExecStart=/opt/openavc/venv/bin/python -m server.main
Restart=always
RestartSec=5
Environment=OPENAVC_DATA_DIR=/var/lib/openavc
Environment=OPENAVC_LOG_DIR=/var/log/openavc
Environment=OPENAVC_PROJECT=/var/lib/openavc/projects/default/project.avc
Environment=OPENAVC_BIND=0.0.0.0
Environment=OPENAVC_ALLOW_ANONYMOUS=false

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
# driver_repo and plugin_repo live under /var/lib/openavc now (covered by
# the data-dir entry). The legacy /opt/openavc/{driver,plugin}_repo paths
# are kept writable with the `-` prefix (tolerated when absent) so the
# runtime migration in server/system_config.migrate_legacy_repos can drain
# them on first start after upgrade from the pre-data_dir layout.
ReadWritePaths=/var/lib/openavc /var/log/openavc -/opt/openavc/driver_repo -/opt/openavc/plugin_repo
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
