# installer/openavc.service
[Unit]
Description=OpenAVC Room Control Server
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=openavc
Group=openavc
WorkingDirectory=/opt/openavc
ExecStartPre=-+/opt/openavc/update-helper.sh /var/lib/openavc
ExecStart=/opt/openavc/venv/bin/python -m server.main
Restart=always
RestartSec=5
Environment=OPENAVC_DATA_DIR=/var/lib/openavc
Environment=OPENAVC_LOG_DIR=/var/log/openavc
Environment=OPENAVC_PROJECT=/var/lib/openavc/projects/default/project.avc
Environment=OPENAVC_BIND=0.0.0.0
# Secure by default: bound to all interfaces so panels can reach it, but the
# Programmer/admin surface requires a credential. A fresh install is
# "unclaimed" — the first Programmer visit sets the admin password. The room
# panel stays open. Set to "true" only if another auth layer fronts OpenAVC.
Environment=OPENAVC_ALLOW_ANONYMOUS=false
# Tells main.py that systemd will relaunch us on exit (Restart=always),
# so cloud-restart can just exit cleanly instead of also spawning its own
# replacement. systemd already sets INVOCATION_ID, but the explicit flag
# is the canonical contract.
Environment=OPENAVC_SERVICE_MANAGED=1

# Security hardening
NoNewPrivileges=true
# Device discovery's ICMP ping sweep has to open an ICMP/raw socket. The
# service runs as the unprivileged `openavc` user, and NoNewPrivileges=true
# strips the file capability from /bin/ping at exec — so the only way the
# sweep can send echo requests is an ambient CAP_NET_RAW granted directly to
# this process (scoped to OpenAVC alone). Without it, the kernel's default
# ping_group_range excludes the service's gid and scans silently return zero
# live hosts.
AmbientCapabilities=CAP_NET_RAW
ProtectSystem=strict
# driver_repo and plugin_repo live under /var/lib/openavc now (covered by
# the data-dir entry). The legacy /opt/openavc/{driver,plugin}_repo paths
# are kept writable with the `-` prefix (tolerated when absent) so the
# runtime migration in server/system_config.migrate_legacy_repos can drain
# them on first start after upgrade from the pre-data_dir layout. Once
# every install has migrated these can be dropped.
ReadWritePaths=/var/lib/openavc /var/log/openavc -/opt/openavc/driver_repo -/opt/openavc/plugin_repo
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
